# luci-theme-argon-jitu/root/etc/uci-defaults/30_luci-theme-argon-jitu
#!/bin/sh
uci batch <<-EOF
	set luci.themes.Argon_jitu=/luci-static/argon_jitu
	commit luci
EOF
exit 0
