# luci-theme-tiktok/root/etc/uci-defaults/30_luci-theme-tiktok
#!/bin/sh
uci batch <<-EOF
	set luci.themes.Tiktok=/luci-static/tiktok
	commit luci
EOF
exit 0
