# luci-theme-tiktok/root/etc/uci-defaults/30_luci-theme-tiktok
#!/bin/sh
if [ "$PKG_UPGRADE" != 1 ]; then
uci batch <<-EOF
	set luci.themes.Tiktok=/luci-static/tiktok
	set luci.main.mediaurlbase=/luci-static/tiktok
	commit luci
EOF
fi
exit 0
